# hub/infra/server/systemd/barranquia-ngrok.service
[Unit]
Description=BarranquIA Hub - Ngrok Tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=desarrollo
ExecStart=/home/desarrollo/.local/bin/ngrok start barranquia --config /home/desarrollo/.config/ngrok/ngrok.yml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
